# debian/gooroom-security-utils.preinst
#!/bin/sh
set -e

/bin/systemctl stop gooroom-security-logant || true
/bin/systemctl disable gooroom-security-logant || true

exit 0
